# hyundai-server.service
[Unit]
Description=Hyundai/Kia Connect API Server (Python)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
# Define the working directory
WorkingDirectory=/home/pi/hyundai-python-server

# Make sure this path to python inside your venv is correct
ExecStart=/home/pi/hyundai-python-server/venv/bin/python3 /home/pi/hyundai-python-server/hyundai_server.py

# Use the full, verified path to the .env file
EnvironmentFile=/home/pi/hyundai-python-server/.env

# Restart policy (comments on separate lines)
Restart=on-failure
# Wait 10 seconds before restarting
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target